# node/supervisord.conf
[program:node]
command = node /home/dotcloud/current/twitter-poll.js
